# overlays/common/usr/lib/systemd/user/commonarch-update-check.service
[Unit]
Description=System updates

[Service]
ExecStart=/usr/bin/system update-check

[Install]
WantedBy=default.target
